# scripts/rootfs-services/copaw.service
[Unit]
Description=CoPaw App
After=network.target

[Service]
ExecStart=/home/tenbox/.copaw/venv/bin/copaw app --host 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
